# amway/lynx-provision/ansible-ua/roles/github-actions/autotest-agent/templates/vnc-workaround.service.j2
[Unit]
After=vncserver@.service
Description=VNC Server Restart

[Service]
Type=oneshot
User=root
ExecStart=/opt/scripts/vnc.sh

[Install]
WantedBy=multi-user.target
