# amway/lynx-provision/ansible/roles/bamboo-agent/templates/bamboo/init_bamboo_agent.service.j2
[Unit]
Description=Daemon to init bamboo agent
After=syslog.target

[Service]
User=root
ExecStart=/opt/scripts/init_bamboo_agent.sh

[Install]
WantedBy=multi-user.target
